# systemd/interview.service
[Unit]
Description=AI Interview Platform Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=root
Group=root
WorkingDirectory=/root/qtale/service/interview/backend
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Load environment variables from .env file
EnvironmentFile=/root/qtale/service/interview/backend/.env

# Gunicorn command
ExecStart=/root/.local/bin/uv run gunicorn \
    --workers 4 \
    --bind unix:/root/qtale/service/interview/backend/interview.sock \
    --timeout 120 \
    --max-requests 1000 \
    --max-requests-jitter 50 \
    --access-logfile /var/log/interview/access.log \
    --error-logfile /var/log/interview/error.log \
    --log-level info \
    config.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/qtale/service/interview/backend /var/log/interview

[Install]
WantedBy=multi-user.target
